# my_tools/app/led_controller/10_check_led.rules
#
# Openvox hardware led udev rules
#
#HW version : V2.0
#

#Create LCD symbol link
DEVPATH=="*devices/pci0000:00/0000:00:1d.0/usb1/1-1/1-1.2/1-1.2.4/1-1.2.4.4/1-1.2.4.4:1.0*", SYMLINK+="opvx/lcd/0", RUN+="/my_tools/check_led.sh"

